# systemd/error-alert.timer
[Unit]
Description=Knowledge System Error Alert Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
