# systemd/evolution.timer
[Unit]
Description=Run Evolution daily at 2:00 AM MSK (23:00 UTC prev day)

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
