# systemd/staleness-check.service
[Unit]
Description=Staleness Check for Knowledge System
After=network.target qdrant.service memgraph.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/.openclaw/workspace/knowledge-system
ExecStart=/usr/bin/node /root/.openclaw/workspace/knowledge-system/src/evolution/staleness-checker.js
StandardOutput=append:/var/log/knowledge/staleness-check.log
StandardError=append:/var/log/knowledge/staleness-check.log

[Install]
WantedBy=multi-user.target
